# moodle-bot.service
[Unit]
Description=Moodle Student Tracker Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/moodle-bot
ExecStart=/opt/moodle-bot/venv/bin/python telegram_bot.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
