# scripts/moodle-bot.service
[Unit]
Description=Moodle Student Tracker Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/moodle-bot
EnvironmentFile=/opt/moodle-bot/.env
ExecStart=/opt/moodle-bot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moodle-bot

[Install]
WantedBy=multi-user.target
